# dockerfiles/final/etc/supervisord-appium.conf
[supervisord]
nodaemon=true

[program:appium]
environment=HOME="/root",DISPLAY=":0",USER="root"
command=/usr/bin/xfce4-terminal -x "/usr/bin/appium"
autorestart=true
priority=5000
