# dockerfiles/final/etc/supervisord-emulator.conf
[supervisord]
nodaemon=true

[program:emulator]
environment=HOME="/root",DISPLAY=":0",USER="root",PWD="/root"
command=/usr/bin/bash -c /etc/launch_emulator.sh
autorestart=true
priority=5000
