# dockerfiles/99_finialize/etc/supervisord-appium.conf
[supervisord]
nodaemon=true

[program:appium]
environment=HOME="/root",DISPLAY=":0",USER="root"
command=/usr/bin/xfce4-terminal -x "bash"
autorestart=true
priority=100
